# dots/bkg.timer
[Unit]
Description=Background Switch Timer

[Timer]
OnCalendar=hourly
Persistent=false
Unit=bkg.service

[Install]
WantedBy=timers.target
